# picam-webcam_1.0-1_all/DEBIAN/postinst
#!/bin/sh
set -e
if [ "$1" = "configure" ]; then
    systemctl daemon-reload
    systemctl enable webcam.service
    systemctl restart webcam.service
fi
exit 0

