# picam-webcam_1.0-1_all/DEBIAN/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop webcam.service || true
    systemctl disable webcam.service || true
fi
exit 0

